# deploy/bot.service
[Unit]
Description=Telegram Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-bot
ExecStart=/home/ubuntu/trading-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment file with secrets
EnvironmentFile=/home/ubuntu/trading-bot/.env

# Kill gracefully on stop
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
